# provisioning/roles/app/templates/npm-dev.service.j2
[Service]
ExecStart=/usr/bin/npm run watch
ExecStop=/bin/kill -INT ${MAINPID}
WorkingDirectory={{ app_path }}
